# roles/docker_host/templates/podman-api-tls.service.j2
[Unit]
Description=Podman Docker-compatible API over mTLS (Jenkins Docker Cloud)
# dbus.service (not just dbus.socket): containers created through this API use
# crun's systemd cgroup manager, which talks to the message bus directly --
# socket-activation alone left a narrow window (seen in testing as "crun:
# cannot open sd-bus") where the very first container creation raced the
# broker's own startup.
After=network-online.target dbus.service
Wants=network-online.target dbus.service

[Service]
ExecStart=/usr/bin/podman system service --time=0 \
  --tls-cert={{ docker_host_certs_dir }}/server-cert.pem \
  --tls-key={{ docker_host_certs_dir }}/server-key.pem \
  --tls-client-ca={{ docker_host_certs_dir }}/ca.pem \
  tcp://0.0.0.0:{{ docker_host_tls_port }}
Restart=on-failure

[Install]
WantedBy=multi-user.target
